# extras/autobuilder.crontab
40 5,11,17,23 * * * /home/oe/gumstix/gumstix-oe/extras/autobuilder.sh > /home/oe/oe_build.log

